# com.matthew.meeting-followup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.matthew.meeting-followup</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>YOUR_HOME_DIR/.meeting-followup/meeting_followup.py</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>ANTHROPIC_API_KEY</key>
        <string>YOUR_ANTHROPIC_API_KEY</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>

    <!-- Watch the Granola data directory (not a specific versioned file).
         Fires when any cache-v*.json is written, so Granola upgrades
         (v3 → v4 → v5) don't require plist changes. -->
    <key>WatchPaths</key>
    <array>
        <string>YOUR_HOME_DIR/Library/Application Support/Granola</string>
    </array>

    <!-- Don't fire more than once every 30 seconds -->
    <key>ThrottleInterval</key>
    <integer>30</integer>

    <key>StandardOutPath</key>
    <string>YOUR_HOME_DIR/.meeting-followup/launchagent-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>YOUR_HOME_DIR/.meeting-followup/launchagent-stderr.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
